# app-utils/mouseemu/autobuild/overrides/usr/lib/systemd/system/mouseemu.service
[Unit]
Description=Mouseemu Initialization

[Service]
Type=simple
ExecStart=/usr/bin/mouseemu

[Install]
WantedBy=multi-user.target
