# lang-js/nvm/autobuild/overrides/etc/profile.d/30-nvm.sh
export NVM_DIR="${HOME}/.nvm"
. /usr/share/nvm/nvm.sh
