# app-web/thunderbird/autobuild/overrides/etc/profile.d/thunderbird.sh
# Whether to use smooth scrolling by default.
#
# Disabled due to a tooltip display bug.
#
# Ref: https://bugzilla.mozilla.org/show_bug.cgi?id=1805703
export MOZ_USE_XINPUT2=1
