# app-i18n/ibus/02-ibus/overrides/etc/xdg/autostart/ibus-daemon.desktop
[Desktop Entry]
Name=IBus daemon startup
Comment=Initialize and daemonize the IBus input framework
Exec=ibus-daemon -d -x
Terminal=false
Type=Application
Categories=
GenericName=
NoShowIn=GNOME;
