# app-utils/elfutils/autobuild/overrides/etc/profile.d/debuginfod.sh
export DEBUGINFOD_URLS="https://debuginfod.elfutils.org"
export DEBUGINFOD_PROGRESS=1
