# extra-devel/emscripten/autobuild/overrides/etc/profile.d/emscripten.sh
export EMSCRIPTEN="/usr/lib/emscripten"
export EMSCRIPTEN_FASTCOMP="/usr/lib/emscripten-fastcomp"

export PATH="$PATH:$EMSCRIPTEN"

