# meta-bases/rocm-base/autobuild/overrides/etc/profile.d/rocm.sh
export ROCM_PATH="/usr/lib/rocm"
export PATH="$PATH:/usr/lib/rocm/bin"
