# app-utils/fzf/autobuild/overrides/etc/profile.d/fzf-key-bindings.bash
#! /usr/bin/bash
source /usr/share/fzf/key-bindings.bash
