# extra-utils/qtchooser/autobuild/overrides/etc/profile.d/qtchooser.sh
case ":${PATH:-}:" in
    *:/usr/lib/qtchooser:*) ;;
    *) PATH="/usr/lib/qtchooser${PATH:+:$PATH}" ;;
esac
