# app-admin/apt/autobuild/overrides/usr/lib/systemd/system/apt-update.service
[Unit]
Description=Check for package update
Documentation=man:apt(8)
ConditionACPower=true

[Service]
Type=oneshot
ExecStart=/usr/bin/apt update
