# app-devel/hook-uname/autobuild/overrides/etc/profile.d/hook_uname.sh
export LD_PRELOAD="$LD_PRELOAD:/usr/lib/hook_uname.so"
