# extra-admin/finger/autobuild/overrides/usr/lib/systemd/system/finger.socket
[Unit]
Description=Finger Socket for Per-Connection Servers

[Socket]
ListenStream=79
Accept=yes

[Install]
WantedBy=sockets.target
