# app-network/samba/autobuild/overrides/usr/lib/systemd/system/smbd@.service
[Unit]
Description=Samba SMB/CIFS server instance

[Service]
ExecStart=/usr/bin/smbd -F
ExecReload=/usr/bin/kill -HUP $MAINPID
StandardInput=socket
