# lang-java/openjdk-11/autobuild/overrides/etc/profile.d/openjdk-11.sh
# FIXME: This have been changed to use a patch in other JDK versions
#        See commit 38a1699a47cd67139db41b7c1d559e345f153166 at github.com/AOSC-Tracking/jdk
#        However the patch is not backported to JDK 11 yet.
export _JAVA_OPTIONS="-Dawt.useSystemAAFontSettings=on"
