# extra-java/openjfx-11/autobuild/overrides/etc/profile.d/openjfx.sh
export _JAVA_OPTIONS="$_JAVA_OPTIONS --module-path=/usr/share/openjfx/lib/"
