# app-admin/trousers/autobuild/overrides/usr/lib/systemd/system/tcsd.service
[Unit]
Description=TCG Core Services Daemon

[Service]
Type=forking
ExecStart=/usr/bin/tcsd

[Install]
WantedBy=multi-user.target
