# extra-libs/usbmuxd/autobuild/overrides/usr/lib/systemd/system/usbmuxd.service
[Unit]
Description=iOS USB Muxer

[Service]
ExecStart=/usr/bin/usbmuxd -u -f -U usbmux
