# extra-utils/topmenu-gtk/autobuild/overrides/etc/profile.d/topmenu-gtk.sh
export GTK_MODULES+=" topmenu-gtk-module "
