# extra-java/gradle/autobuild/overrides/etc/profile.d/gradle.sh
#!/usr/bin/env sh
export GRADLE_HOME=/usr/lib/gradle
export PATH=$PATH:$GRADLE_HOME/bin

