# lang-dotnet/dotnet-host/autobuild/overrides/etc/profile.d/dotnet-cli-telemetry-optout.sh
# https://learn.microsoft.com/dotnet/core/tools/telemetry
export DOTNET_CLI_TELEMETRY_OPTOUT=1
