# extra-java/openjdk-8/autobuild/overrides/etc/profile.d/openjdk-8.sh
export JAVA_HOME="/usr/lib/java-8"
export _JAVA_OPTIONS="-Dawt.useSystemAAFontSettings=on"
export JAVA_FONTS="/usr/share/fonts"
