# app-emulation/latx/autobuild/overrides/usr/lib/systemd/system/usr-gnemul-lat\x2dx86_64-dev.mount
[Unit]
Description=Mount /dev for the LATX runtime

[Mount]
What=/dev
Where=/usr/gnemul/lat-x86_64/dev
Type=none
Options=bind
LazyUnmount=true

[Install]
WantedBy=local-fs.target
