# app-emulation/latx/autobuild/overrides/usr/lib/systemd/system/usr-gnemul-latx\x2dx86_64-run.mount
[Unit]
Description=Mount /run for the LATX runtime

[Mount]
What=/run
Where=/usr/gnemul/latx-x86_64/run
Type=none
Options=bind
LazyUnmount=true

[Install]
WantedBy=local-fs.target
