# extra-multimedia/minitube/autobuild/overrides/etc/profile.d/minitube.sh
export GOOGLE_API_KEY="AIzaSyBGpe01okUSW2GQDIgLovY23Mj1RKzzOOY"
