# app-network/openssh/autobuild/overrides/usr/lib/systemd/system/sshd.socket
[Unit]
Description=OpenSSH Daemon Server Socket
Before=sockets.target
Wants=sshdgenkeys.service

[Socket]
ListenStream=22
Accept=no

[Install]
WantedBy=sockets.target

