# extra-productivity/libreoffice/autobuild/overrides/etc/profile.d/libreoffice.sh
# Uncomment to force a certain visual style in LibreOffice...
#     - gen: Generic/X11
#     - gtk: GTK+ 2.0
#     - gtk3: GTK+ 3.0

#export SAL_USE_VCLPLUGIN=gen
#export SAL_USE_VCLPLUGIN=gtk
#export SAL_USE_VCLPLUGIN=gtk3
