# extra-productivity/libreoffice/autobuild/overrides/etc/profile.d/libreoffice.sh
# Uncomment desired user interface library.

# export SAL_USE_VCLPLUGIN=gen
# export SAL_USE_VCLPLUGIN=gtk
export SAL_USE_VCLPLUGIN=gtk3
