# extra-network/babeld/autobuild/overrides/usr/lib/systemd/system/babeld.service
[Unit]
Description=babeld routing daemon
After=network.target

[Service]
ExecStart=/usr/bin/babeld

[Install]
WantedBy=multi-user.target
