# extra-utils/fzf/autobuild/overrides/etc/profile.d/fzf-key-bindings.zsh
#! /usr/bin/zsh
source /usr/share/fzf/key-bindings.zsh
