# extra-java/openjdk/autobuild/overrides/etc/profile.d/openjdk.sh
export JAVA_HOME="/usr/lib/java"
export _JAVA_OPTIONS="-Dawt.useSystemAAFontSettings=on"
export JAVA_FONTS="/usr/share/fonts"
