# extra-cross/binutils+cross-bm-armel/autobuild/overrides/etc/profile.d/cross-bm-armel.sh
# set PATH for armel bare metal cross toolchain

PATH="$PATH:/opt/bmcross/armel/bin/"
