# app-web/firefox/autobuild/overrides/etc/profile.d/firefox.sh
# Whether to use smooth scrolling by default.
export MOZ_USE_XINPUT2=1
