# app-network/bind/autobuild/overrides/usr/lib/systemd/system/bind-generate-rndc-key.service
[Unit]
Description=Generate rndc key for BIND (named)

[Service]
Type=oneshot
ExecStart=/usr/bin/bind-generate-rndc-key
