# app-network/bind/autobuild/overrides/usr/lib/systemd/system/bind-update-root-hints.timer
[Unit]
Description=Update DNS root hints for BIND (named)

[Timer]
OnCalendar=weekly
OnActiveSec=0
Persistent=true

[Install]
WantedBy=timers.target
