# extra-devel/arduino-mk/autobuild/overrides/etc/profile.d/90-arduino-mk
export AVR_TOOLS_DIR=/
export ARDMK_DIR=/usr/share/arduino
