# app-admin/finger/autobuild/overrides/usr/lib/systemd/system/finger@.service
[Unit]
Description=Finger Per-Connection Server

[Service]
User=nobody
ExecStart=/usr/sbin/in.fingerd
StandardInput=socket
