# extra-x11/devilspie/autobuild/overrides/etc/xdg/autostart/devilspie.desktop
[Desktop Entry]
Name=devilspie
Exec=devilspie
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
