# extra-devel/cuda/autobuild/overrides/etc/profile.d/cuda.sh
export CUDA_PATH="/usr/lib/cuda"
export PATH="$PATH:/usr/lib/cuda/bin"

