# app-web/exim/autobuild/overrides/usr/lib/systemd/system/exim.socket
[Unit]
Description=Exim MTA
Conflicts=exim.service

[Socket]
ListenStream=25
Accept=yes

[Install]
WantedBy=sockets.target

