# app-shells/powershell/autobuild/overrides/etc/profile.d/powershell-telemetry-optout.sh
# https://learn.microsoft.com/powershell/module/microsoft.powershell.core/about/about_telemetry
export POWERSHELL_TELEMETRY_OPTOUT=1
