# app-network/mihomo/autobuild/overrides/usr/lib/systemd/user/clash-meta.service
[Unit]
Description=Starting Mihomo (Clash Meta)
After=network.target

[Service]
Type=exec
Restart=on-abort
ExecStart=/usr/bin/mihomo

[Install]
WantedBy=default.target
