# app-network/vsftpd/autobuild/overrides/usr/lib/systemd/system/vsftpd@.service
[Unit]
Description=vsftpd per-connection FTP server

[Service]
ExecStart=-/usr/bin/vsftpd
ExecReload=/usr/bin/kill -HUP $MAINPID
StandardInput=socket
StandardError=syslog
