# app-network/vsftpd/autobuild/overrides/usr/lib/systemd/system/vsftpd.service
[Unit]
Description=vsftpd FTP daemon
After=network.target

[Service]
ExecStart=/usr/bin/vsftpd
ExecReload=/usr/bin/kill -HUP $MAINPID
KillMode=process

[Install]
WantedBy=multi-user.target
