# app-multimedia/openh264/autobuild/overrides/etc/profile.d/openh264-gmp.sh
# Allow Mozilla applications to find system OpenH264.
export MOZ_GMP_PATH="${MOZ_GMP_PATH}${MOZ_GMP_PATH:+:}/usr/lib/mozilla/plugins/gmp-gmpopenh264/system-installed"
