# app-web/httpd/autobuild/overrides/usr/lib/systemd/system/httpd.socket
[Unit]
Description=Apache HTTP Server Socket

[Socket]
ListenStream=80
NoDelay=true
DeferAcceptSec=30

[Install]
WantedBy=sockets.target
