# extra-network/polipo/autobuild/overrides/usr/lib/systemd/system/polipo.service
[Unit]
Description=Polipo proxy daemon
After=network.target

[Service]
ExecStart=/usr/bin/polipo
User=polipo
PrivateDevices=yes

[Install]
WantedBy=multi-user.target
