# extra-admin/apparmor/autobuild/overrides/usr/lib/systemd/system/apparmor.service
[Unit]
Description=AppArmor Security Daemon
DefaultDependencies=no
Before=sysinit.target

[Service]
Type=oneshot
ExecStart=/usr/bin/apparmor_load
ExecStop=/usr/bin/apparmor_unload
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
