# lang-java/openjdk-common/autobuild/overrides/etc/profile.d/openjdk.sh
export JAVA_HOME="/usr/lib/java"
export JAVA_FONTS="/usr/share/fonts"
