# extra-libs/speech-dispatcher/autobuild/overrides/usr/lib/systemd/system/speech-dispatcherd.service
[Unit]
Description=Speech-Dispatcher daemon
After=syslog.target

[Service]
Type=forking
ExecStart=/usr/bin/speech-dispatcher -d

[Install]
WantedBy=multi-user.target
