# extra-security/trousers/autobuild/overrides/usr/lib/systemd/system/tcsd.service
[Unit]
Description=TCG Core Services Daemon

[Service]
Type=forking
ExecStart=/sbin/tcsd

[Install]
WantedBy=multi-user.target
