# runtime-desktop/qt-5/autobuild/overrides/etc/profile.d/qtwebengine.sh
# FIXME: QtWebEngine sandboxing causes a trace/breakpoint trap exception.
# Disabling sandboxing until we figure this out.
export QTWEBENGINE_DISABLE_SANDBOX=1
