# base-utils/man-db/autobuild/overrides/usr/lib/systemd/system/man-db.timer
[Unit]
Description=Daily Man-DB cache generation

[Timer]
OnCalendar=daily
AccuracySec=1d
Persistent=true
